# addons/point_of_sale/tools/posbox/overwrite_after_init/etc/ld.so.preload
/usr/lib/arm-linux-gnueabihf/libarmmem-v8l.so
