# addons/point_of_sale/tools/posbox/overwrite_before_init/etc/ld.so.preload
# the used qemu-kernel doesn't work well with this library, so get rid of it
#/usr/lib/arm-linux-gnueabihf/libcofi_rpi.so
