# addons/point_of_sale/tools/posbox/overwrite_before_init/etc/ld.so.preload
#/usr/lib/arm-linux-gnueabihf/libcofi_rpi.so
