# addons/point_of_sale/tools/posbox/overwrite_before_init/etc/systemd/system/splashscreen.service
[Unit]
Description=Splash screen
DefaultDependencies=no
After=local-fs.target

[Service]
ExecStart=/usr/bin/fbi -d /dev/fb0 --noverbose -a /var/www/splash.png
StandardInput=tty
StandardOutput=tty

[Install]
WantedBy=sysinit.target
